# config/supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid
user=root

[program:vnc]
command=/usr/bin/vncserver :1 -geometry 1920x1080 -depth 24 -localhost no
autostart=true
autorestart=true
priority=10
stdout_logfile=/var/log/supervisor/vnc.log
stderr_logfile=/var/log/supervisor/vnc_error.log

[program:xfce]
command=/usr/bin/startxfce4
autostart=true
autorestart=true
priority=20
environment=DISPLAY=":1"
stdout_logfile=/var/log/supervisor/xfce.log
stderr_logfile=/var/log/supervisor/xfce_error.log

[program:novnc]
command=/usr/share/novnc/utils/launch.sh --vnc localhost:5901 --listen 6901
autostart=true
autorestart=true
priority=30
stdout_logfile=/var/log/supervisor/novnc.log
stderr_logfile=/var/log/supervisor/novnc_error.log

[program:api]
command=python3 -m uvicorn app.api:app --host 0.0.0.0 --port 8001
directory=/app
autostart=true
autorestart=true
priority=40
stdout_logfile=/var/log/supervisor/api.log
stderr_logfile=/var/log/supervisor/api_error.log

[program:web]
command=python3 -m http.server 3001
directory=/app/web
autostart=true
autorestart=true
priority=50
stdout_logfile=/var/log/supervisor/web.log
stderr_logfile=/var/log/supervisor/web_error.log
